# kronos/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:kronos]
command=sudo service kronosd start
